# 参考项目/looplj-axonhub/deploy/axonhub.service
[Unit]
Description=AxonHub AI Gateway
After=network.target

[Service]
Type=simple
User=axonhub
WorkingDirectory=/opt/axonhub
ExecStart=/opt/axonhub/axonhub
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target